# .config/systemd/user/random_wallpaper.service
[Unit]
Description=Change wallpaper randomly every hour

[Service]
ExecStart=/home/rhyme/.scripts/random_wallpaper.sh
