# .config/systemd/user/random_wallpaper.timer
[Unit]
Description=Run random_wallpaper.service every hour

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
